# deployment/systemd/led-signage.service
[Unit]
Description=LED Matrix Signage Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/PixelCast/led-signage
ExecStart=/usr/bin/python3 /opt/PixelCast/led-signage/daemon.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-signage

# Give GPIO time to settle after boot
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
