# deployment/systemd/PixelCast.service
[Unit]
Description=PixelCast - LED Matrix Signage Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/PixelCast
Environment="SIGNAGE_SECRET=CHANGE_THIS_TO_A_RANDOM_SECRET_KEY"

# ── Pre-start: GPIO settle → mount USB → seed defaults ───────────────────
# The leading '-' means the step is allowed to fail (e.g. no USB attached)
# without preventing the daemon from starting.
ExecStartPre=/bin/sleep 3
ExecStartPre=-/opt/PixelCast/led-signage/deployment/scripts/mount-usb.sh
ExecStartPre=-/opt/PixelCast/led-signage/deployment/scripts/init-usb-data.sh

# ── Daemon ────────────────────────────────────────────────────────────────
# Config + media live on the USB stick (/media/usb).
# If USB is not present the daemon falls back to built-in defaults.
ExecStart=/usr/bin/python3 /opt/PixelCast/led-signage/daemon.py \
    --config=/media/usb/config/panel.json \
    --playlist=/media/usb/config/playlist.json \
    --schedule=/media/usb/config/schedule.json \
    --users=/media/usb/config/users.json \
    --media-dir=/media/usb/media

# Keep Python/Flask on cores 0-2; core 3 (isolcpus=3) stays free for
# the C++ rgbmatrix refresh thread which pins itself there via os.sched_setaffinity
CPUAffinity=0-2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-signage

# Give GPIO time to settle after boot
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
