# deployment/systemd/PixelCast.service
[Unit]
Description=PixelCast - LED Matrix Signage Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/PixelCast
Environment="SIGNAGE_SECRET=CHANGE_THIS_TO_A_RANDOM_SECRET_KEY"
ExecStart=/usr/bin/python3 /opt/PixelCast/led-signage/daemon.py \
    --config=/opt/PixelCast/config/panel.json \
    --playlist=/opt/PixelCast/config/playlist.json \
    --schedule=/opt/PixelCast/config/schedule.json \
    --users=/opt/PixelCast/config/users.json \
    --media-dir=/opt/PixelCast/media
# Keep Python/Flask on cores 0-2; core 3 (isolcpus=3) stays free for
# the C++ rgbmatrix refresh thread which pins itself there via os.sched_setaffinity
CPUAffinity=0-2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-signage

# Give GPIO time to settle after boot
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
